# dist/vpsfree-irc-bot@.service
[Unit]
Description=IRC bot for vpsFree.cz's channels
Wants=network-online.target
After=network-online.target

[Service]
User=vpsfbot
Group=vpsfbot
ExecStart=/opt/vpsfree-irc-bot/bin/vpsfree-irc-bot --config /etc/vpsfree-irc-bot/%i.yml

[Install]
WantedBy=default.target
